# Files/multi-login.service
[Unit]
Description=Dropbear User-Limit Service

[Service]
Type=simple
ExecStart=/usr/local/bin/multi-login

[Install]
WantedBy=multi-user.target